# systems/deimos/rootfs/etc/systemd/system/logid.service.d/override.conf
[Service]
Restart=on-failure
RestartSec=2
StartLimitBurst=10
StartLimitIntervalSec=60
