# systems/cachy/config/systemd/kraken-liquid-fan-curve.service
[Unit]
Description=Apply Kraken liquid fan curve
After=multi-user.target

[Service]
Type=simple
EnvironmentFile=/etc/nix-dotfiles/env
ExecStart=/usr/bin/bash -lc 'exec "$DOTFILES_REPO_ROOT/systems/cachy/bin/kraken-liquid-fan-curve.sh"'
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
